# dot_config/systemd/user/break-reminder.service
[Unit]
Description=Break Reminder Service
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/share/chezmoi/scripts/executable_break-reminder.sh auto
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-1

[Install]
WantedBy=default.target